# sshws/ws-dropbear.service
[Unit]
Description=SSH Over Websocket (Dropbear)
Documentation=https://givps.com
After=network.target nss-lookup.target

[Service]
# Peningkatan Keamanan: Port 2095 adalah high port, jadi jalankan service sebagai user 'nobody'
# untuk membatasi hak akses. Ini lebih aman daripada menggunakan 'root'.
Type=simple
User=nobody
WorkingDirectory=/usr/local/bin
# Menghapus 'CapabilityBoundingSet' dan 'AmbientCapabilities' karena tidak diperlukan oleh user 'nobody'
# saat menggunakan high port (>1024).

# PERBAIKAN KRITIS: Port harus diteruskan menggunakan flag '-p'.
# Ditambahkan juga flag '-b 0.0.0.0' untuk eksplisit binding address.
ExecStart=/usr/bin/python3 -O /usr/local/bin/ws-dropbear -b 0.0.0.0 -p 2095
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5s

# Resource limits
LimitNOFILE=50000
LimitNPROC=5000

[Install]
WantedBy=multi-user.target