# sshws/ws-stunnel.service
[Unit]
Description=SSH Over Websocket (Stunnel)
Documentation=https://givps.com
After=network.target nss-lookup.target

[Service]
# Peningkatan Keamanan: Port 8443 adalah high port, jadi jalankan service sebagai user 'nobody'
Type=simple
User=nobody
WorkingDirectory=/usr/local/bin

# PERBAIKAN KRITIS: Port harus diteruskan menggunakan flag '-p'. 
# Menggunakan port 8443 sebagai contoh, Anda bisa menggantinya jika port Anda berbeda.
ExecStart=/usr/bin/python3 -O /usr/local/bin/ws-stunnel -b 0.0.0.0 -p 8443
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5s

# Resource limits
LimitNOFILE=50000
LimitNPROC=5000

[Install]
WantedBy=multi-user.target